# django_project/danadjango/deploy/telegram-metapko-bot.service
[Unit]
Description=Telegram bot (Meta-PKO / Qazaq Dana API client)
After=network.target meta-university.service

[Service]
Type=simple
User=meta-university
Group=www-data
WorkingDirectory=/var/www/meta-university
EnvironmentFile=-/var/www/meta-university/.env.telegram
ExecStart=/var/www/meta-university/venv/bin/python /var/www/meta-university/scripts/telegram_metapko_bot.py
Restart=always
RestartSec=8

[Install]
WantedBy=multi-user.target
